# crontab
0 * * * * cd /path/to/your/project && python manage.py get_rosters

# An empty line is required at the end of this file for a valid cron file.
